# configs/crontab
@reboot  /home/pi/piBoxBot/bash/startup.sh &
