# configs/crontab
@reboot /home/pi/piBoxBot/bash/startup.sh &
@reboot /home/pi/piBoxBot/bash/cameraLoop.sh &
